# src/configs/supervisord.conf
[supervisord]
logfile=/tmp/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB       ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10          ; (num of main logfile rotation backups;default 10)
loglevel=info               ; (log level;default info; others: debug,warn,trace)
pidfile=/tmp/supervisord.pid; (supervisord pidfile;default supervisord.pid)
nodaemon=true              ; (start in foreground if true;default false)
minfds=1024                 ; (min. avail startup file descriptors;default 1024)
minprocs=200                ; (min. avail process descriptors;default 200)

# [program:ros_launch]
# command=roslaunch src/configs/perception.launch
# stdout_logfile=/dev/stdout
# stdout_logfile_maxbytes=0
# autorestart=true

[program:roscore]
command=roscore
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
autorestart=true

[program:realsense_capture]
command=python -m src.robots.realsense_image_capture
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
autorestart=true


# [program:pilots_eye]
# command=python -m src.utils.pilots_eye.pilots_eye --port=9091
# stdout_logfile=/dev/stdout
# stdout_logfile_maxbytes=0
# autorestart=true
